# supervisord.conf
[supervisord]

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:mysqld]
command = /usr/local/bin/docker-entrypoint.sh mysqld
directory = /usr/src/app
user = mysql
stdout_logfile = /usr/src/app/mysqld_stdout.log
stderr_logfile = /usr/src/app/mysqld_stderr.log
environment = MYSQL_ROOT_PASSWORD=root, MYSQL_USER=db_user, MYSQL_PASSWORD=db_pass, MYSQL_DATABASE=distance

[program:app]
command = python3 -m swagger_server
directory = /usr/src/app
user = root
stdout_logfile = /usr/src/app/app_stdout.log
stderr_logfile = /usr/src/app/app_stderr.log
environment = LC_ALL=C.UTF-8, LANG=C.UTF-8, FLASK_APP=app.py, MYSQL_USER=db_user, MYSQL_PASSWORD=db_pass, MYSQL_DATABASE=distance, PRODUCTION=1
