# supervisord.conf
[supervisord]

[unix_http_server]
file=/usr/src/app/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///usr/src/app/supervisor.sock

[program:neo4j]
command = /sbin/tini -g -s -- /docker-entrypoint.sh neo4j
directory = /usr/src/app
stdout_logfile = /usr/src/app/neo4j_stdout.log
stderr_logfile = /usr/src/app/neo4j_stderr.log
environment = NEO4J_EDITION=community, NEO4J_AUTH=none
priority = 1

[program:app]
command = python3 -m swagger_server
directory = /usr/src/app
stdout_logfile = /usr/src/app/app_stdout.log
stderr_logfile = /usr/src/app/app_stderr.log
environment = DEBUG=1
priority = 999
